# assets/system/auto-suspend-service
#!/sbin/openrc-run
# Este archivo no se linkea, se instala con install -m 755
# Al hacer cambios no olvides volverlo a instalar con:
# sudo install -m 755 "$HOME/.dotfiles/assets/system/auto-suspend-service" /etc/init.d/auto-suspend

description="Suspender el sistema cuando la batería está baja"
command="/usr/local/bin/auto-suspend-loop"
command_args=""
pidfile="/var/run/auto-suspend.pid"
start_stop_daemon_args="-b -m --pidfile ${pidfile}"

depend() {
	need localmount
	after bootmisc
}
